# rimay-type.service
[Unit]
Description=rimay-type speech-to-keyboard
After=graphical-session.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
ExecStart=%h/.cargo/bin/rimay-type
Restart=on-failure
RestartSec=3
PassEnvironment=DISPLAY XAUTHORITY

[Install]
WantedBy=default.target
